# install/debian12/mail_api.service
[Unit]
Description=MAIL_API (webhook receiver + control panel)
After=network.target

[Service]
Type=simple
User=mail_api
Group=mail_api
WorkingDirectory=/opt/mail_api/app
Environment=MAIL_API_DATA_DIR=/var/lib/mail_api
ExecStart=/opt/mail_api/venv/bin/python -m mail_api
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
